# build/package/debian/blitzd.service
# RaspiBlitz Info Daemon: systemd unit for blitzd
# /etc/systemd/system/blitzd.service

[Unit]
Description=RaspiBlitz Daemon
After=network.target

[Service]
User=bitcoin
Group=bitcoin
ExecStart=/usr/local/bin/blitzd

[Install]
WantedBy=multi-user.target